# ops/com.aicompany.dispatcher-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.aicompany.dispatcher-poller</string>
  <key>ProgramArguments</key>
  <array><string>/bin/bash</string><string>/Users/nurikim/.aicompany/dispatcher-poller.sh</string></array>
  <key>StartInterval</key><integer>60</integer>
  <key>RunAtLoad</key><true/>
  <key>LimitLoadToSessionType</key><string>Aqua</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>LANG</key><string>en_US.UTF-8</string>
  </dict>
  <key>WorkingDirectory</key><string>/Users/nurikim</string>
  <key>StandardOutPath</key><string>/Users/nurikim/Library/Logs/aicompany/poller-launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/nurikim/Library/Logs/aicompany/poller-launchd.log</string>
  <key>ProcessType</key><string>Standard</string>
</dict>
</plist>
